# systemd/tg-monitor.service
[Unit]
Description=Telegram Server Monitor (Node Exporter -> Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/tg-monitoring
EnvironmentFile=%h/tg-monitoring/.env
# Default: use new modular entrypoint (tgbot.main)
ExecStart=%h/tg-monitoring/.venv/bin/python -m tgbot.main
# Legacy fallback: ExecStart=%h/tg-monitoring/.venv/bin/python -m monitor.main
Restart=always
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=default.target
